# etcs/alpha/supervisor.conf
[program:public-git]
command=node app.js
numprocs=1
directory=/home/nodejs/public-git/public-git
autostart=true
autorestart=true
stopsignal=QUIT
stopwaitsecs=120
user=nodejs
